# autostart/awm.desktop
[Desktop Entry]
Type=Application
Terminal=false
Exec=avant-window-navigator --startup &
Name=Avant-Window-Navigator
