# autostart/conky.desktop
[Desktop Entry]
Encoding=UTF-8
Type=Application
Name=conky
Exec=conky
Terminal=false
