# autostart/redshift-gtk.desktop
[Desktop Entry]
GenericName[lt]=Spalvos temperatūros reguliavimas
Name[pl]=Redshift
GenericName[nl]=Kleurtemperatuur aanpassing
Name[zh_CN]=红移
Comment[zh_CN]=色温调节工具
Name[ru]=Redshift
Name[de]=Redshift
Name[pt_BR]=Redshift
Comment[tr]=Renk sıcaklığı ayarlama aracı
Type=Application
Comment[es]=Herramienta para el ajuste de la temperatura de color
Comment[da]=Justeringsværktøj for farvetemperatur
Comment[fr]=Outil de réglage de la température de couleur
Comment[ru]=Инструмент регулирования цветовой температуры
Categories=Utility;
GenericName[es]=Ajuste de la temperatura de color
GenericName[fr]=Réglage de la température de couleur
X-GNOME-Autostart-enabled=false
Exec=redshift-gtk
Name[it]=Redshift
GenericName[zh_CN]=色温调节
Hidden=true
Name[fr]=Redshift
Comment[lt]=Spalvos temperatūros reguliavimo įrankis
Comment[nl]=Kleurtemperatuur aanpassingsprogramma
Name=Redshift
GenericName[cs]=Přizpůsobení barevné teploty
Version=1.0
GenericName[de]=Farbtemperatur Anpassung
Comment[de]=Farbtemperatur Anpassungswerkzeug
Name[es]=Redshift
GenericName[da]=Justering af farvetemperatur
Icon=redshift
GenericName[pl]=Dostosowanie temperatury barwowej
Comment=Color temperature adjustment tool
GenericName[hu]=Szín hőmérséklet beállítás
GenericName[tr]=Renk sıcaklığı ayarı
Name[hu]=Redshift
Comment[pl]=Narzędzie do dostosowywania temperatury barwowej
Name[tr]=Redshift
GenericName=Color temperature adjustment
Comment[hu]=Színhőmérséklet beállító eszköz
Comment[cs]=Nástroj pro přizpůsobení barevné teploty
StartupNotify=true
GenericName[it]=Regolazione della temperatura del colore
Terminal=false
Name[lt]=Redshift
GenericName[ru]=Настройка цветовой температуры
Comment[it]=Strumento per la regolazione della temperatura del colore

